# svn-server/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:svnserve]
command=svnserve -d --foreground -r /home/svn/repo
user=svn
stderr_logfile=/home/svn/stderr.txt
